# srcs/influxdb/srcs/supervisord.conf
[supervisord]
nodaemon = true

[program:influxdb]
command = sh -c "influxd && kill -s SIGTERM $(cat supervisord.pid)"

[program:telegraf]
command = sh -c "telegraf --config /etc/telegraf/telegraf.conf && kill -s SIGTERM $(cat supervisord.pid)"