# srcs/ftps/srcs/supervisord.conf
[supervisord]
nodaemon=true

[program:vsftpd]
# command= sh -c "vsftpd /etc/vsftpd.conf && kill -s SIGTERM $(cat supervisord.pid)"
command= sh -c "vsftpd /etc/vsftpd.conf"
autostart=true
autorestart=true