# ftps/srcs/supervisord.conf
[supervisord]
nodaemon=true

[program:nginx]
command= sh -c "vsftpf /etc/vsftpd.conf"
autostart=true
autorestart=true