# wordpress/srcs/supervisord.conf
# [supervisord]
# nodaemon=true

# [program:php-fpm7]
# command=sh -c "php-fpm7"
# autostart=true
# autorestart=true

# [program:nginx]
# command=sh -c "nginx -g 'daemon off;'"
# autostart=true
# autorestart=true

[supervisord]
nodaemon=true
user = root

[program:nginx]
command= sh -c "nginx -g 'daemon off;' && kill -s SIGTERM $(cat supervisord.pid)"
# command= sh -c "nginx -g 'daemon off;'"
# autostart=true
# autorestart=true

[program:php]
command= sh -c "/usr/sbin/php-fpm7 -F && kill -s SIGTERM $(cat supervisord.pid)"
# command= sh -c "/usr/sbin/php-fpm7 -F"
# autostart=true
# autorestart=true
