# wordpress/srcs/supervisord.conf
[supervisord]
nodaemon=true

[program:mysql]
command=sh -c "sleep infinite"
autostart=true
autorestart=true

# [program:telegraf]
# command=sh -c "telegraf"
# autostart=true
# autorestart=true