# grafana old/supervisord.conf
[supervisord]
nodaemon = true
[program:grafana]
command = sh -c "grafana-server --homepath /grafana-7.5.4 --config /grafana-7.5.4/conf/grafana.ini && kill -s SIGTERM $(cat supervisord.pid)"